# home/.config/systemd/user/astra-commit.timer
[Unit]
Description=Nightly astra-config auto-commit + off-box push (floor; SessionStart hook commits the rest)

[Timer]
OnCalendar=*-*-* 03:00:00 America/Los_Angeles
Persistent=true

[Install]
WantedBy=timers.target
