# home/.config/systemd/user/astra-commit.timer
[Unit]
Description=Nightly astra-config auto-commit (floor; SessionStart hook does the rest)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
