# home/.config/systemd/user/drift-check.service
[Unit]
Description=Check that zaz-astra still reproduces from the astra-config rebuild

[Service]
Type=oneshot
ExecStart=/root/astra-config/scripts/drift-check.sh /root/astra-config
