# home/.config/systemd/user/drift-check.timer
[Unit]
Description=Daily astra-config drift check (drops a SessionStart banner sentinel on drift)

[Timer]
OnCalendar=*-*-* 05:30:00 America/Los_Angeles
Persistent=true

[Install]
WantedBy=timers.target
