# home/.config/systemd/user/git-access-check.timer
[Unit]
Description=Daily git PAT repo-access probe (12:30 PT)

[Timer]
OnCalendar=*-*-* 12:30:00 America/Los_Angeles
Persistent=true

[Install]
WantedBy=timers.target
