# home/.config/systemd/user/health-check.timer
[Unit]
Description=Hourly grok-mcp funnel liveness check (emails after 2 consecutive failures)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
