# home/.config/systemd/user/health-check.service
[Unit]
Description=Probe the grok-mcp funnel and email on sustained failure

[Service]
Type=oneshot
ExecStart=/root/astra-config/scripts/health-check.sh /root/astra-config
