# home/.config/systemd/user/journal-oauth-watch.service
[Unit]
Description=Watch journaling Claude OAuth; email if broken ≥48h

[Service]
Type=oneshot
ExecStart=/root/astra-config/scripts/journal-oauth-watch.sh
