# home/.config/systemd/user/ops-log.service
[Unit]
Description=Machine-readable ops snapshot (agent-run ops-log)

[Service]
Type=oneshot
ExecStart=/root/astra-config/scripts/agent-run.sh ops-log
